# config/supervisor.d/squid.ini
[program:squid]
command=/usr/sbin/squid -sYCNd 1 -f /etc/squid/squid.conf
priority=900
stdout_logfile=/var/log/supervisord-squid.log
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/supervisord-squid.log
stderr_logfile_maxbytes=0
username=proxy
autorestart=true
